# rootfs/etc/services.d/packetstream/run
#!/command/with-contenv bash
exec /bin/packetstream_service 
sleep 5