# rootfs/etc/services.d/packetstream/run
#!/command/with-contenv sh
/bin/packetstream_service


