# rootfs/etc/services.d/packetstream/run
#!/command/with-contenv bash
su picash -c '/bin/packetstream_service'


