# rootfs/etc/services.d/bitping/run
#!/command/with-contenv bash
exec /bin/bitping_service
sleep 5