# rootfs/etc/services.d/pawnsapp/run
#!/command/with-contenv bash
exec /bin/pawnsapp_service 
sleep 5