# rootfs/etc/services.d/pawnsapp/run
#!/command/with-contenv bash
su picash -c "/bin/pawnsapp_service"
sleep 5