# rootfs/etc/services.d/honeygain/run
#!/command/with-contenv sh
/bin/honeygain_service


