# rootfs/etc/services.d/honeygain/run
#!/command/with-contenv bash
exec /bin/honeygain_service


