# rootfs/etc/services.d/traffmonetizer/run
#!/command/with-contenv bash
su picash -c "/bin/traffmonetizer_service"


