# rootfs/etc/services.d/traffmonetizer/run
#!/command/with-contenv sh
/bin/traffmonetizer_service


