# rootfs/etc/services.d/traffmonetizer/run
#!/command/with-contenv bash
exec /bin/traffmonetizer_service 
sleep 5