# rootfs/etc/services.d/earnapp/run
#!/command/with-contenv bash
exec /bin/earnapp_service 
sleep 5