# rootfs/etc/services.d/peer2profit/run
#!/command/with-contenv bash
su picash -c "/bin/peer2profit_service"


