# rootfs/etc/services.d/peer2profit/run
#!/command/with-contenv bash
exec /bin/peer2profit_service 
sleep 5
